# root/etc/s6-overlay/s6-rc.d/init-mod-universal-ssmtp/run
#!/usr/bin/with-contenv bash

echo "\
    ssmtp" >> /mod-repo-packages-to-install.list
